# srcpd/10-opendcc.rules
SYSFS{idVendor}=="0403", SYSFS{idProduct}=="bfd8", RUN+="/sbin/modprobe -q ftdi_sio vendor=0x0403 product=0xbfd8"
